# crontab
* * * * * /usr/bin/flock -n /tmp/ytdownload.lock youtube-dl --config-location /config/youtube-dl.conf >> /dev/stdout 2>&1



