# crontab
* * * * * youtube-dl --version >> /dev/stdout 2>&1



